# deploy/systemd/mm-healthcheck@.timer
[Unit]
Description=Run MagicMirror stack health check every minute for user %i

[Timer]
OnBootSec=90
OnUnitActiveSec=60
Unit=mm-healthcheck@%i.service
Persistent=true

[Install]
WantedBy=timers.target
