# deploy/systemd/quran-voice@.service
[Unit]
Description=MMM Quran Voice Listener for user %i
After=network-online.target sound.target ollama.service
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=20

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/MagicMirror-Pi5/magicmirror/modules/MMM-QuranDisplay
Environment=PYTHONUNBUFFERED=1
Environment=OMP_NUM_THREADS=2
Environment=OPENBLAS_NUM_THREADS=2
Environment=NUMEXPR_NUM_THREADS=2
Environment=HEARTBEAT_FILE=/home/%i/MagicMirror-Pi5/magicmirror/modules/MMM-QuranDisplay/logs/voice_listener.heartbeat
Environment=VOICE_DEVICE=pulse
Environment=VOICE_SOURCE=alsa_input.usb-ME6S_MS_N-B_R-UN_ME6S-00.mono-fallback
ExecStartPre=/usr/bin/test -x /home/%i/MagicMirror-Pi5/magicmirror/modules/MMM-QuranDisplay/start_listener.sh
ExecStart=/bin/bash /home/%i/MagicMirror-Pi5/magicmirror/modules/MMM-QuranDisplay/start_listener.sh
Restart=always
RestartSec=2
TimeoutStartSec=45
TimeoutStopSec=15
KillSignal=SIGTERM
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=6
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
